# myarchy/.config/systemd/user/feh-wallpaper.timer
[Unit]
Description=Random wallpaper with feh

[Timer]
OnCalendar=hourly
Unit=feh-wallpaper.service

[Install]
WantedBy=default.target
